# robot-twoeyes/services/crestereo.service
[Unit]
Description=CREStereo Disparity Estimation API Server
After=network.target

[Service]
Type=simple
User=unitree
WorkingDirectory=/home/unitree/Documents/yx/project/robot-twoeyes
Environment="PATH=/home/unitree/miniconda3/envs/yolo/bin:/usr/local/bin:/usr/bin:/bin"
Environment="LD_LIBRARY_PATH=/home/unitree/miniconda3/envs/yolo/lib"
ExecStart=/home/unitree/miniconda3/envs/yolo/bin/python img_process/crestereo/run_server.py --model img_process/crestereo/model/crestereo_combined_iter10_480x640.onnx
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
